# supervision/systemd/visio-watchdog.service
[Unit]
# C10 — watchdog de la chaîne d'ingestion (go2rtc -> ingress -> livekit).
# Démarre APRÈS la stack visio. Préfixe visio-* (jamais les units d'une autre intégration).
Description=videocalltizen — watchdog supervision (C10)
Documentation=file:///home/marco/videocalltizen/supervision/RUNBOOK.md
After=visio-stack.service network-online.target
Wants=visio-stack.service network-online.target

[Service]
Type=simple
User=marco
WorkingDirectory=/home/marco/videocalltizen/supervision
# URLs hôte par défaut (la stack publie 1984/9080/7880 sur localhost).
Environment=GO2RTC_URL=http://localhost:1984
Environment=TOKEN_SERVICE_URL=http://localhost:9080
Environment=LIVEKIT_URL=http://localhost:7880
Environment=ROOM_NAME=salon
Environment=CAMERA_IDENTITY=camera-salon
Environment=GO2RTC_STREAM=salon
Environment=STATUS_PORT=9095
ExecStart=/usr/bin/python3 /home/marco/videocalltizen/supervision/watchdog.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
